# devops/files/nginx/certbot-renew.timer
[Unit]
Description=Renews letsencrypt certs

[Timer]
OnUnitActiveSec=10m

[Install]
WantedBy=timers.target
