# devops/files/nginx/certbot-renew.timer
[Unit]
Description=Renews letsencrypt certs

[Timer]
OnUnitActiveSec=12h 17m

[Install]
WantedBy=timers.target
