# devops/files/nginx/letsencrypt.service
[Unit]
Description=Renews letsencrypt certs
After=syslog.target network.target

[Service]
ExecStart=/usr/bin/certbot renew --quiet
