# devops/packer/nginx/letsencrypt.service
[Unit]
Description=Renews letsencrypt certs
After=syslog.target network.target

[Service]
ExecStart=certbot renew --quiet
